# files/supervisor.d/super-sys.ini
[program:php-fpm]
command=/usr/local/php72/sbin/php-fpm --fpm-config /usr/local/php72/etc/php-fpm.conf
stdout_logfile=/var/log/supervisor-php-fpm.log
redirect_stderr=true
process_name=%(program_name)s_%(process_num)02d
autostart=true
autorestart=true
numprocs=1
user=root
startsecs=1
startretries=100000

[program:nginx]
command=/usr/local/nginx/sbin/nginx -g "daemon off;" -c /usr/local/nginx/conf/nginx.conf
stdout_logfile=/var/log/supervisor-nginx.log
redirect_stderr=true
process_name=%(program_name)s_%(process_num)02d
autostart=true
autorestart=true
numprocs=1
user=root
startsecs=1
startretries=100000

[program:redis]
command=redis-server
stdout_logfile=/var/log/supervisor-redis.log
redirect_stderr=true
process_name=%(program_name)s_%(process_num)02d
autostart=true
autorestart=true
numprocs=1
user=root
startsecs=1
startretries=100000

