# files/supervisor.d/super-queue.ini
[program:bloated-queue-default]
command=/data/httpd/espier-bloated/artisan doctrine:queue:work --queue=default --delay=3 --memory=128 --timeout=30 --sleep=1 --tries=3
stdout_logfile=/var/log/supervisor-bloated-queue-default.log
redirect_stderr=true
process_name=%(program_name)s_%(process_num)02d
autostart=true
autorestart=true
numprocs=1
user=www
startsecs=1
startretries=100000

[program:bloated-queue-seckill]
command=/data/httpd/espier-bloated/artisan doctrine:queue:work --queue=seckill --delay=3 --memory=128 --timeout=30 --sleep=1 --tries=3
stdout_logfile=/var/log/supervisor-bloated-queue-seckill.log
redirect_stderr=true
process_name=%(program_name)s_%(process_num)02d
autostart=true
autorestart=true
numprocs=1
user=www
startsecs=1
startretries=100000

[program:bloated-queue-slow]
command=/data/httpd/espier-bloated/artisan doctrine:queue:work --queue=slow --delay=3 --memory=128 --timeout=1800 --sleep=1 --tries=3
stdout_logfile=/var/log/supervisor-bloated-queue-slow.log
redirect_stderr=true
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=true
autostart=true
autorestart=true
numprocs=1
user=www
startsecs=1
startretries=100000

[program:bloated-queue-sms]
command=/data/httpd/espier-bloated/artisan doctrine:queue:work --queue=sms --delay=3 --memory=128 --timeout=1800 --sleep=1 --tries=3
stdout_logfile=/var/log/supervisor-bloated-queue-sms.log
redirect_stderr=true
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=true
autostart=true
autorestart=true
numprocs=1
user=www
startsecs=1
startretries=100000
